# ha-macmini-poweron/rootfs/etc/services.d/poweron/finish
#!/command/with-contenv sh
#
# poweron finish script — cleanup when the service exits.
#

echo "[$(date '+%Y-%m-%d %H:%M:%S')] [poweron] [addon] finishing..." >&1
